# roles/environment/default_visual_text_editor/templates/etc/profile.d/00.visual.sh.j2
#!/bin/sh

export VISUAL={{ default_visual_text_editor_executable_path }}
