# roles/environment/default_cli_text_editor/templates/etc/profile.d/00.editor.sh.j2
#!/bin/sh

export EDITOR={{ default_cli_text_editor_executables_map[default_cli_text_editor] }}
